# deploy/systemd/xiaozhi-web.service
[Unit]
Description=xiaozhi-bridge web admin (static file server)
Documentation=https://github.com/yourusername/xiaozhi-bridge
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/projects/xiaozhi-bridge/web
# V1: serve via Python http.server (simple, no extra deps)
# V2: use Caddy / nginx instead
ExecStart=/usr/bin/python3 -m http.server 3000 --bind 127.0.0.1 --directory /root/projects/xiaozhi-bridge/web/dist
Restart=always
RestartSec=5
StandardOutput=append:/var/log/xiaozhi-bridge/web.log
StandardError=append:/var/log/xiaozhi-bridge/web.err

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log/xiaozhi-bridge

[Install]
WantedBy=multi-user.target
